# config/octeontx/etc/inittab
::sysinit:/etc/startup.sh

::shutdown:/etc/shutdown.sh

tty1::askfirst:/bin/sh
tty2::askfirst:/bin/sh
tty3::askfirst:/bin/sh
tty4::askfirst:/bin/sh

# Login prompt on @CONSOLE@:
#::respawn:/sbin/getty -L @CONSOLE@ 115200 vt100

# Open root shell on @CONSOLE@:
#vanishing#@CONSOLE@::askfirst:/bin/ash --login
